# .config/systemd/user/backup.timer
[Unit]
Description=Run Backup hourly
Requires=backup.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
