# .config/systemd/user/backup.service
[Unit]
Description=Backup projects
Wants=backup.timer

[Service]
ExecStart=/home/fabio/bin/backup-projects.sh

[Install]
WantedBy=multi-user.target
